# abf-repoclosure/repoclosure.timer
[Unit]
Description=Run motd.service every day at 08.00

[Timer]
OnCalendar=*-*-* 05:00:00

[Install]
WantedBy=timers.target
