# abf-repoclosure/repoclosure.service
[Unit]
Description=Repoclosure report

[Service]
Type=oneshot
ExecStart=/usr/bin/docker run -d -v /home/omv/repoclosure:/repoclosure-report openmandriva/repoclosure
